# etc/s6-overlay/s6-rc.d/05-ssnc/run
#!/command/with-contenv sh
echo "Starting shairport-sync"
exec /usr/local/bin/shairport-sync
