# etc/s6-overlay/s6-rc.d/02-avahi-daemon/run
#!/command/with-contenv sh
echo "Starting avahi-daemon"
exec /usr/sbin/avahi-daemon -D
